# install/systemd/honeypy@.service
[Unit]
Description=HoneyPy - Instância %i
Documentation=https://github.com/yourusername/honeypy
After=network.target
Wants=network.target

[Service]
Type=simple
User=honeypy-%i
Group=honeypy-%i
RuntimeDirectory=honeypy-%i
RuntimeDirectoryMode=0750

Environment=PYTHONUNBUFFERED=1
Environment=HONEYPY_CONFIG=/etc/honeypy/config-%i.json

WorkingDirectory=/opt/honeypy
ExecStart=/usr/local/bin/honeypy --monitor --config /etc/honeypy/config-%i.json
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=honeypy-%i

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/honeypy-%i /var/lib/honeypy-%i
ReadOnlyPaths=/etc/honeypy
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target